# meta-evb/meta-evb-aspeed/meta-evb-ast2600/recipes-phosphor/snmp/files/snmp.service
[Unit]
Description=SNMP
After=snmpd.service

[Service]
Type=oneshot
ExecStart=/usr/bin/snmp.sh start

[Install]
WantedBy=snmpd.service
