# meta-evb/meta-evb-aspeed/meta-evb-ast2600/recipes-phosphor/power-sqlites/files/total_consumption.service
[Unit]
Description=Total Consumption Service
After=network.target

[Service]
Type=simple
ExecStart=/usr/bin/total_consumption
Restart=always

[Install]
WantedBy=multi-user.target